# crontab
# m h  dom mon dow   command
# run at every minute
* * * * * php $HOME/_projects/zz_work/transfermate/index.php >/dev/null 2>&1
